# bin/crontab
#minute hour    day     month   week    command
0       0       *       *       *       /usr/local/bin/backup